# supervisord.conf
[supervisord]
nodaemon=true

[program:dse]
command=service dse start
